# etc/service/postgresql/run
#!/bin/bash -e

exec chpst -u postgres:postgres:ssl-cert /usr/lib/postgresql/9.1/bin/postgres -D /var/lib/postgresql/9.1/main -c config_file=/etc/postgresql/9.1/main/postgresql.conf 2>&1
